# debian/tee-supplicant.service
[Unit]
Description=TEE Supplicant

[Service]
Type=simple
ExecStart=/usr/sbin/tee-supplicant

[Install]
WantedBy=multi-user.target
